# consumers/agent_actions_consumer.service
[Unit]
Description=Agent Actions Kafka Consumer - PostgreSQL Persistence
Documentation=file:///opt/omniclaude/consumers/README.md
After=network.target kafka.service postgresql.service
Wants=kafka.service postgresql.service

[Service]
Type=simple
User=omniclaude
Group=omniclaude
WorkingDirectory=/opt/omniclaude/consumers

# Python environment (adjust path to your virtualenv)
Environment="PATH=/opt/omniclaude/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/opt/omniclaude"

# Kafka configuration
Environment="KAFKA_BROKERS=localhost:9092"
Environment="KAFKA_GROUP_ID=agent-actions-postgres"

# PostgreSQL configuration
Environment="POSTGRES_HOST=localhost"
Environment="POSTGRES_PORT=5436"
Environment="POSTGRES_DATABASE=omniclaude"
Environment="POSTGRES_USER=postgres"
# POSTGRES_PASSWORD must be set externally — use EnvironmentFile or systemctl set-environment
# Do NOT hardcode credentials here; this file is tracked in git

# Consumer tuning
Environment="BATCH_SIZE=100"
Environment="BATCH_TIMEOUT_MS=1000"

# Monitoring
Environment="HEALTH_CHECK_PORT=8080"
Environment="LOG_LEVEL=INFO"

# Start consumer
ExecStart=/opt/omniclaude/venv/bin/python /opt/omniclaude/consumers/agent_actions_consumer.py

# Restart policy
Restart=on-failure
RestartSec=10s
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=100%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-actions-consumer

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/omniclaude

[Install]
WantedBy=multi-user.target
